# daemon/templates/systemd.service
[Unit]
Description=MMU VPN Tray Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart={{EXEC_PATH}}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
