# daemon/mmuvpn.service
[Unit]
Description=MMU VPN Tray Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/mmuvpn
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
